# everything-presence-mmwave-configurator/services/backend/run
#!/command/with-contenv bashio

exec python3 /app/backend.py