# everything-presence-mmwave-configurator/rootfs/etc/s6-overlay/s6-rc.d/zone-configurator/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

cd /app
export NODE_ENV="${NODE_ENV:-production}"
APP_PORT_CONFIG="$(bashio::config 'port')"
if [ -n "${APP_PORT_CONFIG}" ] && [ "${APP_PORT_CONFIG}" != "null" ]; then
  export PORT="${APP_PORT_CONFIG}"
else
  export PORT="${PORT:-3000}"
fi
FIRMWARE_LAN_PORT_CONFIG="$(bashio::config 'firmware_lan_port')"
if [ -n "${FIRMWARE_LAN_PORT_CONFIG}" ] && [ "${FIRMWARE_LAN_PORT_CONFIG}" != "null" ]; then
  export FIRMWARE_LAN_PORT="${FIRMWARE_LAN_PORT_CONFIG}"
else
  export FIRMWARE_LAN_PORT="${FIRMWARE_LAN_PORT:-38080}"
fi

bashio::log.info "Starting Zone Configurator backend on port ${PORT}"
if [ "${PORT}" != "3000" ]; then
  bashio::log.warning "Ingress is configured for port 3000. Changing the app port may disable ingress access."
fi
bashio::log.info "LAN Firmware Server will start on port ${FIRMWARE_LAN_PORT}"

exec node backend/dist/index.js
