# everything-presence-mmwave-configurator/rootfs/etc/s6-overlay/s6-rc.d/zone-configurator/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

cd /app
export NODE_ENV="${NODE_ENV:-production}"
export PORT="${PORT:-3000}"

bashio::log.info "Starting Zone Configurator backend on port ${PORT}"

exec node backend/dist/index.js
